# deploy/systemd/titop-arenda-backup.service
[Unit]
Description=Titop Arenda PostgreSQL backup
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/flock --wait 300 /run/lock/titop-arenda-postgres-maintenance.lock /usr/local/sbin/titop-arenda-postgres-backup
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
