# deploy/systemd/titop-arenda-backup-alert@.service
[Unit]
Description=Report Titop Arenda backup failure to GlitchTip
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/titop-arenda-glitchtip-alert %i
